# fl.conf
[program:fl]
command=python3 /home/curator4/fl/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/fl/test.err.log
stdout_logfile=/home/curator4/fl/test.out.log